# supervisord example configuration.conf
; This file is part of the Quake Live server implementation by TomTec Solutions. Do not copy or redistribute or link to this file without the emailed consent of Thomas Jones (thomas@tomtecsolutions.com).
; Quake Live supervisord configuration file

[unix_http_server]           ; local socket file
file=/tmp/supervisor.sock    ; (the path to the socket file)

;[inet_http_server]          ; inet (TCP) server disabled by default
;port=*:9001                 ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; socket file
;serverurl=http://127.0.0.1:9001      ; inet server

[program:quakelive]
command=/home/qlserver/quakestart.sh %(process_num)s           ; Sent to 'quakestart.sh', from 0 to numprocs, see what happens by looking in 'quakestart.sh'. 
user=qlserver                                                  ; Default Quake Live server account.
process_name=qzeroded_%(process_num)s                          ; Name the Quake Live dedicated server process, 'qzeroded_numproc'.
numprocs=13                                                    ; How many server processes should be spawned.
autorestart=true                                               ; If the servers crash, or the process exits, start them up again.

[program:quakelive_serverbuilder]
command=/home/qlserver/quakestart_serverbuilder.sh %(process_num)s
user=qlserver
process_name=qzeroded_%(process_num)s
numprocs=12
autorestart=true
